# infra/deploy/systemd/sentinel-alpha.service
[Unit]
Description=Sentinel Alpha API
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/sentinel-alpha
EnvironmentFile=/opt/sentinel-alpha/.env
ExecStart=/usr/bin/python3 -m uvicorn apps.api.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
